# src/zm-metric.service.in
# This is a skeleton created by zproject.
# You can add hand-written code here.

[Unit]
Description=zm-metric-snmp service
Requires=network.target malamute.service
After=network.target malamute.service
Before=zm-asset.service
PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@prefix@/share/bios/etc/default/bios
EnvironmentFile=-@prefix@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@prefix@/share/zm/etc/default/zm
EnvironmentFile=-@prefix@/share/zm/etc/default/zm__%n.conf
EnvironmentFile=-@sysconfdir@/default/bios
EnvironmentFile=-@sysconfdir@/default/bios__%n.conf
EnvironmentFile=-@sysconfdir@/default/zm
EnvironmentFile=-@sysconfdir@/default/zm__%n.conf
Environment="prefix=@prefix@"
ExecStart=@prefix@/bin/zm-metric-snmp --rules '/var/lib/zm/zm-metric-snmp/rules'

[Install]
WantedBy=bios.target
